# photoapp.conf
description "photoapp"

start on started mountall

stop on shudown

respawn 
respawn limit 99 5


exec node /home/ubuntu/photoapp/app.js >> /var/log/photoapp.log 2>&1




